# roles/fxo/files/fxo.timer
[Unit]
Description=check and clear FXO lines periodically

[Timer]
OnActiveSec=7200

[Install]
WantedBy=timers.target
