# roles/firewall/files/phonewall.timer
[Unit]
Description=phonewall closer

[Timer]
Unit=phonewall_close.service
OnActiveSec=1h

[Install]
WantedBy=multi-user.target
